# bedtime-reminder.service
[Unit]
Description=Bedtime Reminder Notification
After=sound.target

[Service]
Type=oneshot
User=%u
Environment="DISPLAY=:0"
ExecStart=/home/%u/.local/bin/bedtime-reminder.sh
