# bedtime-reminder.timer
[Unit]
Description=Bedtime Reminder Timer

[Timer]
# Run at 11:00 PM every day
OnCalendar=*-*-* 23:00:00
Persistent=false

[Install]
WantedBy=timers.target
